# init/policyd-rate-limit-clean.timer
[Unit]
Description=Clean policyd rate limit database daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
