# debian/policyd-rate-limit.cron.daily
#!/bin/sh
# cron script for policyd-rate-limit
test -x /usr/bin/policyd-rate-limit || exit 0
/usr/bin/policyd-rate-limit --clean >/dev/null
